# contrib/mail-admin-v2-reputation.timer
[Unit]
Description=Mail Admin v2 — Reputation Snapshot (hourly)

[Timer]
OnCalendar=hourly
Persistent=true
Unit=mail-admin-v2-reputation.service

[Install]
WantedBy=timers.target
