# contrib/mail-admin-v2-reputation.service
[Unit]
Description=Mail Admin v2 — Reputation Snapshot (one-shot)
After=mail-admin-v2.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mail-admin-v2-reputation-snapshot
