# systemd/mail-admin-v2-dnsbl-snapshot.timer
[Unit]
Description=Run DNSBL outbound snapshot every 6 hours

[Timer]
OnBootSec=10min
OnUnitActiveSec=6h
AccuracySec=2min
Persistent=true
Unit=mail-admin-v2-dnsbl-snapshot.service

[Install]
WantedBy=timers.target
